# deploy/supervisor/project-x-worker.conf
; Laravel queue worker for project-x (test-slack backend).
; Install on the server:  sudo cp this file to /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start project-x-worker:*
; Without a running worker, ALL queued jobs (notifications, push, Ai translation,
; reel video processing) silently never run when QUEUE_CONNECTION is redis/database.

[program:project-x-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/project-x/backend/artisan queue:work --queue=default,heavy1,heavy2,heavy3,sendComment,increment-diamond,luckyBox,lucky_gift,lucky_gift_2,lucky_gift_3 --sleep=3 --tries=3 --max-time=3600 --backoff=10
directory=/var/www/project-x/backend
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/project-x/backend/storage/logs/worker.log
stopwaitsecs=3600
